# hk_web_app/rootfs/etc/services.d/hk-web-app/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# HK Web App – Node-Server (Ingress + WebSocket-Proxy zum HA-Core)
# ==============================================================================
exec node /app/server.mjs
