# hk_addon/rootfs/etc/services.d/hk-addon/run
#!/usr/bin/with-contenv sh
# ==============================================================================
# HK Addon – Node-Server (Ingress + Proxy zum HA-Core)
# sh + with-contenv: Supervisor-Env wie bashio, ohne Bash-Pflicht für reines exec
# ==============================================================================
# Intl/ICU (Europe/Berlin) braucht genug V8-Heap; auf manchen Add-on-Hosts sonst
# „DateTimePatternGeneratorCache … out of memory“ → 502 durch ständigen Neustart.
export NODE_OPTIONS="--max-old-space-size=128"
exec node /app/server.mjs
