# hk_addon/rootfs/etc/services.d/hk-addon/run
#!/usr/bin/with-contenv sh
# ==============================================================================
# HK Addon – Node-Server (Ingress + Proxy zum HA-Core)
# sh + with-contenv: Supervisor-Env wie bashio, ohne Bash-Pflicht für reines exec
# ==============================================================================
exec node /app/server.mjs
