# hk_addon/rootfs/etc/services.d/hk-addon/run
#!/usr/bin/with-contenv sh
# ==============================================================================
# HK Addon – Node-Server (Ingress + Proxy zum HA-Core)
# sh + with-contenv: Supervisor-Env wie bashio, ohne Bash-Pflicht für reines exec
# ==============================================================================
# Europe/Berlin ohne Intl: Date#getHours() etc. nutzen TZ (tzdata im Image).
export TZ=Europe/Berlin
export NODE_OPTIONS="--max-old-space-size=128"
exec node /app/server.mjs
