# hk_addon/rootfs/etc/services.d/hk-addon/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# HK Addon – Node-Server (Ingress + Proxy zum HA-Core)
# ==============================================================================
exec node /app/server.mjs
